# statuslist.service
[Unit]
Description=Service to start StatusList

[Service]
Type=simple
ExecStart=/home/pi/StatusList/run.py

[Install]
WantedBy=multi-user.target
